# conf/jellyfin-backup.service
[Unit]
Description="Jellyfin Backup"

[Service]
Type=oneshot
ExecStart=/usr/bin/jellyman -ba

