# conf/jellyfin-backup.service
[Unit]
Description="Jellyfin Backup"

[Service]
ExecStart=/usr/bin/jellyman -ba
